# deploy/yocto/meta-loopy/recipes-loopy/loopy-bundle/files/loopy.service
# Runs the Loopy app as a Wayland client after weston is up. weston-init starts
# weston (kiosk-shell + the HDMI outputs from weston.ini); this unit connects to
# its Wayland socket.
#
# NOTE (blocked-verify): the Wayland socket name (wayland-0/wayland-1), the
# XDG_RUNTIME_DIR, and the user weston runs as depend on the weston-init setup on
# the image — confirm on device and adjust WAYLAND_DISPLAY / XDG_RUNTIME_DIR / User.
[Unit]
Description=Loopy floor-console app (Wayland client under weston)
After=weston.service
Requires=weston.service

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=WAYLAND_DISPLAY=wayland-1
ExecStart=/usr/bin/loopy-kiosk-launch
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
