# deploy/yocto/meta-loopy/recipes-loopy/loopy-bundle/files/loopy-wireplumber.service
# WirePlumber session manager for the Loopy appliance (see loopy-pipewire.service).
[Unit]
Description=WirePlumber (Loopy appliance)
After=loopy-pipewire.service
Requires=loopy-pipewire.service
[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/wireplumber
Restart=always
RestartSec=2
[Install]
WantedBy=multi-user.target
