# deploy/yocto/meta-loopy/recipes-loopy/loopy-bundle/files/loopy.service
# Runs the Loopy app as a Wayland client after weston is up. weston-init starts
# weston (kiosk-shell + the HDMI outputs from weston.ini); this unit connects to
# its Wayland socket.
#
# NOTE (blocked-verify): the Wayland socket name (wayland-0/wayland-1), the
# XDG_RUNTIME_DIR, and the user weston runs as depend on the weston-init setup on
# the image — confirm on device and adjust WAYLAND_DISPLAY / XDG_RUNTIME_DIR / User.
[Unit]
Description=Loopy floor-console app (Wayland client under weston)
After=weston.service loopy-wireplumber.service
Requires=weston.service
# Audio: order after PipeWire/WirePlumber so the launcher's pw-jack finds the
# server. Wants (not Requires) — the app still runs (silent) if audio is down.
Wants=loopy-wireplumber.service

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-1
ExecStart=/usr/bin/loopy-kiosk-launch
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
